# packaging/darwin/postinstall
#! /bin/sh

date >/tmp/post

sudo /bin/launchctl load /Library/LaunchDaemons/org.zeek.zeek-agent.plist

exit 0
